# packaging/systemd/agent-secret-bridge.service
[Unit]
Description=Agent Secret Bridge isolated credential broker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=asb
Group=asb-clients
UMask=0007
EnvironmentFile=/etc/agent-secret-bridge/asb.env
ExecStart=/usr/local/bin/asb daemon --config /etc/agent-secret-bridge/asb.yaml --config-sha256 ${ASB_CONFIG_SHA256} --socket /run/agent-secret-bridge/broker.sock
Restart=on-failure
RestartSec=2s

RuntimeDirectory=agent-secret-bridge
RuntimeDirectoryMode=0750
StateDirectory=agent-secret-bridge
StateDirectoryMode=0700

NoNewPrivileges=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
ReadWritePaths=/run/agent-secret-bridge /var/lib/agent-secret-bridge
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
